# sentinel/southport-sentinel.service
[Unit]
Description=Southport Sentinel Service
After=network.target

[Service]
ExecStart=/usr/bin/node /home/pi/southport-system/sentinel/index.js

[Install]
WantedBy=multi-user.target
